# novnc-supervisord.conf

[supervisord]
nodaemon=true

[program:X11]
command=Xvfb :0 -screen 0 425x900x24
autorestart=true

[program:x11vnc]
command=/usr/bin/x11vnc
autorestart=true

[program:novnc]
command=/usr/src/noVNC/utils/launch.sh --vnc localhost:5900 --listen 8080
autorestart=true

[program:socat]
command=socat tcp-listen:6000,reuseaddr,fork unix:/tmp/.X11-unix/X0
autorestart=true

[program:mirai]
command=java -jar ./mcl.jar
directory=/usr/src/app
autorestart=true
